# kpi-engine-app.service
[Unit]
Description="kpi-engine-app"
After=network.target

[Service]
User=root
WorkingDirectory=/root/kpi-engine/
ExecStart=/usr/bin/python3 app.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
